# systemd/blit@.service
[Unit]
Description=blit terminal multiplexer for %i
Requires=blit@%i.socket

[Service]
Type=simple
User=%i
ExecStart=blit-server
Restart=no
